# node_exporter.conf
# Run node_exporter
# Run prometheus

start on startup

chdir /opt/node_exporter

script
./node_exporter
end script
